# deploy/systemd/stock-analyze-monthly-summary.service
[Unit]
Description=Stock Analyze monthly evolution reminder
After=network-online.target stock-analyze-monthly-review.service
Wants=network-online.target
OnFailure=stock-analyze-pipeline-failure@%n.service

[Service]
Type=oneshot
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
EnvironmentFile=-/etc/stock-analyze/secrets.env
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli --logs-dir /opt/stock-analyze/logs notify-workflow-summary --cadence monthly --repo-root /opt/stock-analyze/app
StandardOutput=append:/opt/stock-analyze/logs/monthly-summary.log
StandardError=append:/opt/stock-analyze/logs/monthly-summary.err
