# deploy/systemd/stock-analyze-qdii-research.timer
[Unit]
Description=Run QDII event and shadow research after weekly strategy jobs

[Timer]
OnCalendar=Sat *-*-* 10:30:00 Asia/Shanghai
Persistent=true
AccuracySec=1min
Unit=stock-analyze-qdii-research.service

[Install]
WantedBy=timers.target

