# deploy/systemd/stock-analyze-codex-weekly.timer
[Unit]
Description=Run Codex agent weekly after Friday close (offset 5 min from Claude)

[Timer]
OnCalendar=Fri *-*-* 09:45:00 UTC
Persistent=true
Unit=stock-analyze-codex-weekly.service

[Install]
WantedBy=timers.target
