# deploy/systemd/stock-analyze-claude-daily.timer
[Unit]
Description=Run Claude agent daily after A-share close

[Timer]
OnCalendar=Mon..Fri *-*-* 09:30:00 UTC
Persistent=true
Unit=stock-analyze-claude-daily.service

[Install]
WantedBy=timers.target
