# deploy/systemd/stock-analyze-claude-cn-qdii-etf-daily.service
[Unit]
Description=Stock Analyze daily simulation run (Defensive mainland QDII ETF)
After=network-online.target
Wants=network-online.target
OnFailure=stock-analyze-pipeline-failure@%n.service
OnSuccess=stock-analyze-aggregate-dashboard.service

[Service]
Type=oneshot
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
EnvironmentFile=-/etc/stock-analyze/secrets.env
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli --market cn_qdii_etf --agent claude run-daily
StandardOutput=append:/opt/stock-analyze/logs/claude-cn-qdii-etf-daily.log
StandardError=append:/opt/stock-analyze/logs/claude-cn-qdii-etf-daily.err
