# deploy/systemd/stock-analyze-claude-daily.service
[Unit]
Description=Stock Analyze daily simulation run (Claude agent)
After=network-online.target
Wants=network-online.target
OnFailure=stock-analyze-pipeline-failure@%n.service
OnSuccess=stock-analyze-aggregate-dashboard.service

[Service]
Type=oneshot
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli --agent claude --logs-dir /opt/stock-analyze/logs run-daily --offline
StandardOutput=append:/opt/stock-analyze/logs/claude-daily.log
StandardError=append:/opt/stock-analyze/logs/claude-daily.err
