# deploy/systemd/stock-analyze-claude-weekly.service
[Unit]
Description=Stock Analyze weekly simulation run (Claude agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli --agent claude --logs-dir /opt/stock-analyze/logs run-weekly
StandardOutput=append:/opt/stock-analyze/logs/claude-weekly.log
StandardError=append:/opt/stock-analyze/logs/claude-weekly.err
