# deploy/systemd/stock-analyze-weekly-trigger.timer
[Unit]
Description=Sat 10:00 Asia/Shanghai weekly-agent dispatch (uses Friday's cache)
# 02:00 UTC == 10:00 Asia/Shanghai (no DST)

[Timer]
OnCalendar=Sat *-*-* 02:00:00 UTC
Persistent=false
AccuracySec=30s
Unit=stock-analyze-weekly-trigger.service

[Install]
WantedBy=timers.target
