# deploy/systemd/stock-analyze-claude-weekly.timer
[Unit]
Description=Run Claude agent weekly after Friday close

[Timer]
OnCalendar=Fri *-*-* 09:40:00 UTC
Persistent=true
Unit=stock-analyze-claude-weekly.service

[Install]
WantedBy=timers.target
