# deploy/systemd/stock-analyze-market-data.timer
[Unit]
Description=Mon-Fri 18:30 Asia/Shanghai market-data fetch + agent dispatch
# 10:30 UTC == 18:30 Asia/Shanghai (no DST). Moved from 17:25 -> 18:30 on
# 2026-05-29: Tushare's same-day EOD daily/daily_basic is not always published
# by 17:25, which made spot() come back empty and fail the fetch (see the
# spot-empty retry in markets/a_share/market_data.py for the in-process guard).

[Timer]
OnCalendar=Mon..Fri *-*-* 10:30:00 UTC
Persistent=false
AccuracySec=30s
Unit=stock-analyze-market-data.service

[Install]
WantedBy=timers.target
