# deploy/systemd/stock-analyze-dashboard.service
[Unit]
Description=Stock Analyze dashboard server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
# serve-dashboard expects ``reports/`` (root reports dir) so it can resolve
# /pro.html → reports/competition/dashboard.html and /simple/<agent>.html
# under reports/competition/simple/. Previous --reports-dir pointing at the
# competition subdir made the path joins double-nest. Omit the flag and let
# WorkingDirectory anchor the lookup at reports/competition/* relatively.
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli serve-dashboard --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/stock-analyze/logs/dashboard.log
StandardError=append:/opt/stock-analyze/logs/dashboard.err

[Install]
WantedBy=multi-user.target
