# deploy/systemd/stock-analyze-codex-weekly.service
[Unit]
Description=Stock Analyze weekly simulation run (Codex agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/stock-analyze/app
Environment=TZ=Asia/Shanghai
ExecStart=/opt/stock-analyze/venv/bin/python -m stock_analyze.cli --agent codex --logs-dir /opt/stock-analyze/logs run-weekly --offline
StandardOutput=append:/opt/stock-analyze/logs/codex-weekly.log
StandardError=append:/opt/stock-analyze/logs/codex-weekly.err
